# config/chainwatch-enforce.service
[Unit]
Description=nullbot under chainwatch seccomp/eBPF enforcement
Documentation=https://github.com/ppiankov/nullbot-dist
After=network-online.target pastewatch-proxy.service
Wants=network-online.target pastewatch-proxy.service
ConditionPathExists=/sys/kernel/btf/vmlinux

[Service]
Type=simple
ExecStart=/usr/local/bin/chainwatch enforce \
    --profile nullbot-enforce \
    --fallback \
    --audit-log /var/log/chainwatch/enforce.jsonl \
    -- /usr/local/bin/nullbot daemon
Restart=on-failure
RestartSec=30
EnvironmentFile=-/root/.workledger/api-key.env
EnvironmentFile=-/root/.nullbot/llm.env
# Route LLM traffic through pastewatch proxy for secret redaction
Environment=NULLBOT_LLM_BASE_URL=http://127.0.0.1:8443

# eBPF requires elevated capabilities
AmbientCapabilities=CAP_BPF CAP_PERFMON CAP_SYS_RESOURCE
CapabilityBoundingSet=CAP_BPF CAP_PERFMON CAP_SYS_RESOURCE CAP_DAC_READ_SEARCH

# Security hardening
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/log/chainwatch /sys/fs/bpf /tmp /root/.nullbot/state
PrivateTmp=true
WorkingDirectory=/tmp

[Install]
WantedBy=multi-user.target
