# config/nullbot.service
[Unit]
Description=nullbot fleet observer
Documentation=https://github.com/ppiankov/nullbot-dist
After=network-online.target pastewatch-proxy.service
Wants=network-online.target pastewatch-proxy.service

[Service]
Type=simple
ExecStart=/usr/local/bin/nullbot daemon
Restart=on-failure
RestartSec=30
EnvironmentFile=-/root/.workledger/api-key.env
EnvironmentFile=-/root/.nullbot/llm.env
# Route LLM traffic through pastewatch proxy for secret redaction
Environment=NULLBOT_LLM_BASE_URL=http://127.0.0.1:8443
WorkingDirectory=/tmp

# Security hardening
NoNewPrivileges=false
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/tmp /var/log/chainwatch /root/.nullbot/state
PrivateTmp=true

[Install]
WantedBy=multi-user.target
