# config/nullbot.service
[Unit]
Description=nullbot fleet observer
Documentation=https://github.com/ppiankov/nullbot-dist
After=network-online.target pastewatch-proxy.service pastewatch-hiveram.service
Wants=network-online.target pastewatch-proxy.service pastewatch-hiveram.service

[Service]
Type=simple
# Wait for proxies to be ready before starting nullbot
ExecStartPre=/bin/sh -c 'for i in 1 2 3 4 5; do curl -s -o /dev/null http://127.0.0.1:8443 2>/dev/null && break || sleep 1; done'
ExecStartPre=/bin/sh -c 'for i in 1 2 3 4 5; do curl -s -o /dev/null http://127.0.0.1:8444 2>/dev/null && break || sleep 1; done'
ExecStart=/usr/local/bin/nullbot daemon
Restart=on-failure
RestartSec=30
EnvironmentFile=-/root/.workledger/api-key.env
EnvironmentFile=-/root/.nullbot/llm.env
EnvironmentFile=-/root/.nullbot/profile.env
# Route ALL outbound traffic through pastewatch proxies for secret redaction
# LLM traffic → :8443, Hiveram traffic → :8444
Environment=NULLBOT_LLM_BASE_URL=http://127.0.0.1:8443
Environment=WORKLEDGER_URL=http://127.0.0.1:8444
WorkingDirectory=/tmp

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/tmp /root/.nullbot/state
PrivateTmp=true

[Install]
WantedBy=multi-user.target
